# deploy/prunebck.service
[Unit]
Description=prune dm database backups

[Service]
Type=oneshot
ExecStart=/dm/bin/prunedaily /dm/prod/bck
